# Misc/service.service
[Unit]
Description=
After=multi-user.target

[Service]
WorkingDirectory=
ExecStart=

[Install]
WantedBy=multi-user.target
